# launchd/com.messagebot.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.messagebot.assistant</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/divagarwal/.bun/bin/bun</string>
      <string>/Users/divagarwal/Projects/message-bot/src/cli.ts</string>
      <string>start</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>BUN_INSTALL</key>
      <string>/Users/divagarwal/.bun</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/divagarwal/Projects/message-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/divagarwal/Projects/message-bot/data/assistant.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/divagarwal/Projects/message-bot/data/assistant.stderr.log</string>
  </dict>
</plist>
